# addons/wordpress/rootfs/etc/services.d/wordpress/run
#!/usr/bin/env bash

# Avvia i servizi
echo "Avvio dei servizi..."
exec s6-svscan /etc/services.d
