# addons/wordpress/rootfs/etc/services.d/wordpress/run
#!/usr/bin/env bash

# Avvia PHP-FPM
php-fpm81

# Avvia Nginx
nginx -g "daemon off;"
