# addons/wordpress/rootfs/etc/services.d/wordpress/run
#!/usr/bin/env bash

# Avvia PHP-FPM in background
php-fpm81

# Avvia Nginx in foreground
exec nginx -g "daemon off;"
